# horseplus/rootfs/etc/s6-overlay/s6-rc.d/horseplus/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the HorsePlus Booking web app
# ==============================================================================

cd /app
exec python3 /app/main.py
